# clickup-mcp/clickup-mcp.service
[Unit]
Description=ClickUp MCP Server (local task cache)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/projects/clickup-mcp
ExecStart=%h/projects/clickup-mcp/venv/bin/python %h/projects/clickup-mcp/main.py
Restart=on-failure
RestartSec=10

StandardOutput=append:%h/.clickup-mcp/sync.log
StandardError=append:%h/.clickup-mcp/sync.log

[Install]
WantedBy=default.target
